# deploy/systemd/gunicorn-entrade.service
[Unit]
Description=Gunicorn — ENTRADE (Django em 127.0.0.1:8000)
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/home/APICLEAR/ENTRADE
Environment="PATH=/home/APICLEAR/ENTRADE/.venv/bin"
ExecStart=/home/APICLEAR/ENTRADE/.venv/bin/gunicorn \
    --config /home/APICLEAR/ENTRADE/gunicorn_config.py \
    ENTRADE.wsgi:application

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
